# bin/2.etcd/templates/etcd.service.j2
[Unit]
Description=Etcd Server
After=network.target
After=network-online.target
Wants=network-online.target
Documentation=https://github.com/coreos

[Service]
Type=notify
WorkingDirectory={{k8s_install_home}}/etcd/
ExecStart=/usr/local/bin/etcd \
 --name={{etcdname}} \
 --data-dir={{k8s_install_home}}/etcd \
 --cert-file=/etc/kubernetes/ssl/etcd.pem \
 --key-file=/etc/kubernetes/ssl/etcd-key.pem \
 --trusted-ca-file=/etc/kubernetes/ssl/ca.pem \
 --peer-cert-file=/etc/kubernetes/ssl/etcd.pem \
 --peer-key-file=/etc/kubernetes/ssl/etcd-key.pem \
 --peer-trusted-ca-file=/etc/kubernetes/ssl/ca.pem \
 --peer-client-cert-auth \
 --client-cert-auth \
 --listen-peer-urls=https://{{ip}}:2380 \
 --initial-advertise-peer-urls=https://{{ip}}:2380 \
 --listen-client-urls=https://{{ip}}:2379,http://127.0.0.1:2379 \
 --advertise-client-urls=https://{{ip}}:2379 \
 --initial-cluster-token=etcd-cluster-0 \
 --initial-cluster={{etcdclustes}} \
 --initial-cluster-state=new
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target